# docker/crontab
35 * * * * php /app/src/feed_parser.php
