# docker/root.cron
35 5 * * * php /app/src/feed_parser.php
*/5 * * * * php /app/src/worker.php
