# docker/root.cron
35 * * * * php /app/src/feed_parser.php
0,30 * * * * php /app/src/worker.php
