# docker/root.cron
35 5 * * * php /app/src/feed_parser.php
* * * * * php /app/src/worker.php
